# data_2_local/sync_stock_list/sync_stock_list.timer
[Unit]
Description=Run sync_stock_list.sh periodically

[Timer]
OnCalendar=*-*-* 17:50:01
Persistent=true

[Install]
WantedBy=timers.target
