# data_2_local/dividend/sync_dividend.timer
[Unit]
Description=Run sync_dividend.sh periodically

[Timer]
# 每周五18:02
OnCalendar=Fri 18:02:00
Persistent=true

[Install]
WantedBy=timers.target
