# data_2_local/dividend/sync_dividend.timer
[Unit]
Description=Run sync_dividend.sh periodically

[Timer]
OnCalendar=Fri 18:22:00
Persistent=true

[Install]
WantedBy=timers.target
