# data_2_local/dividend/sync_stock_fhps_detail_em.timer
[Unit]
Description=Run sync_stock_fhps_detail_em.sh periodically

[Timer]
OnCalendar=Fri 18:21:00
Persistent=true

[Install]
WantedBy=timers.target
