# data_2_local/sh/bfq_daily_etf_price_2_local/sync_bfq_daily_etf_price.timer
[Unit]
Description=Run sync_bfq_daily_etf_price.sh periodically

[Timer]
OnCalendar=*-*-* 18:02:00
Persistent=true

[Install]
WantedBy=timers.target
