# data_2_local/sync_stock_share_change/sync_stock_share_change.timer
[Unit]
Description=Run sync_stock_share_change.sh periodically

[Timer]
OnCalendar=*-*-* 18:25:01
Persistent=true

[Install]
WantedBy=timers.target
